# 69-dm-integrity.rules
# do not edit this file, it will be overwritten on update

# Bring up dm integrity devices

ACTION=="remove", GOTO="dm_integrity_end"

SUBSYSTEM!="block", GOTO="dm_integrity_end"
ENV{ID_FS_TYPE}!="DM_integrity", GOTO="dm_integrity_end"

ENV{ID_PART_ENTRY_UUID}=="?*", RUN+="/usr/sbin/integritysetup_udev $env{DEVNAME} $env{ID_PART_ENTRY_UUID}", GOTO="dm_integrity_end"
ENV{ID_WWN}=="?*", RUN+="/usr/sbin/integritysetup_udev $env{DEVNAME} $env{ID_WWN}", GOTO="dm_integrity_end"

LABEL="dm_integrity_end"
